# deploy/waypoint-vps/waypoint-control.service
[Unit]
Description=Waypoint control API
After=network-online.target tailscaled.service
Wants=network-online.target
Requires=tailscaled.service

[Service]
Type=simple
User=waypoint-mitm
Group=waypoint-mitm
WorkingDirectory=/opt/waypoint
EnvironmentFile=/etc/waypoint/waypoint-control.env
RuntimeDirectory=waypoint
RuntimeDirectoryMode=0750
ExecStart=/opt/waypoint/.venv/bin/python /opt/waypoint/tools/waypoint_control_api.py --host ${WAYPOINT_CONTROL_HOST} --port ${WAYPOINT_CONTROL_PORT} --state-dir ${WAYPOINT_STATE_DIR} --runtime-dir ${WAYPOINT_RUNTIME_DIR}
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
