# deploy/waypoint-vps/waypoint-mitm.service
[Unit]
Description=Waypoint gs-loc transparent MITM
After=network-online.target tailscaled.service
Wants=network-online.target
Requires=tailscaled.service

[Service]
Type=simple
User=waypoint-mitm
Group=waypoint-mitm
WorkingDirectory=/opt/waypoint
EnvironmentFile=/etc/waypoint/waypoint-spoof.env
Environment=HOME=/opt/waypoint
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/waypoint/.venv/bin/mitmdump --mode transparent --listen-port 8080 --showhost --set block_global=false --set flow_detail=0 --allow-hosts=^gs-loc(-cn)?[.]apple[.]com(:443)?$ -s /opt/waypoint/tools/mitm_location_probe.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
